# debian/buienradar.service
[Unit]
Description=Run buienradar
Wants=buienradar.timer

[Service]
Type=oneshot
ExecStart=/usr/bin/buienradar

[Install]
WantedBy=multi-user.target
